# rc.local
# /etc/rc.local

ifconfig 6to4-wan6 mtu 1434

ifconfig eth1 txqueuelen 16

(
	ethtool -K eth1 tso off
	ethtool -K eth1 gso off
	ethtool -K eth1 ufo off
	ethtool -K eth1 gro off
	ethtool -K eth1 lro off
	ethtool -G eth1 tx 8
) 2> /dev/null

for I in /sys/class/net/eth1/queues/tx-*/byte_queue_limits/limit_max
	do
		echo 1514 > $I
	done

/etc/debloat.sh eth0 1000 100 64
/etc/debloat.sh wlan0 450 50 16
/etc/debloat.sh wlan1 450 50 16
/etc/simple-qos.sh

exit 0
